# superv.conf
[program:sensor2mqtt]
command=/home/pi/src/sensormqtt/sensor2mqtt.sh
directory=/home/pi/src/sensormqtt
environment=PATH="/home/pi/src/mqttenv"
user=pi
autostart=true
