# superv_gpio.conf
[program:gpio2mqtt]
command=/home/pi/src/sensormqtt/gpio2mqtt.sh
directory=/home/pi/src/sensormqtt
environment=PATH="/home/pi/src/mqttenv"
user=root
autostart=true
